# node-v8.9.4-alpine/service/mongodb/run
#!/bin/sh

exec /usr/bin/mongod --quiet --bind_ip 0.0.0.0 --dbpath /var/lib/mongodb --logpath /var/log/mongodb/mongod.log
